# pay-bot/deploy/pay-bot.service
[Unit]
Description=pay-bot — relays course payment leads to Telegram for osokin.ai
After=network.target

[Service]
Type=simple
User=pay-bot
Group=pay-bot
WorkingDirectory=/opt/pay-bot
EnvironmentFile=/opt/pay-bot/.env
ExecStart=__NODE__ /opt/pay-bot/server.js
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
